# ty[PWN 150]/ubuntu-arm64/var/lib/dpkg/info/cron.conffiles
/etc/crontab
/etc/cron.weekly/.placeholder
/etc/cron.hourly/.placeholder
/etc/cron.daily/.placeholder
/etc/init/cron.conf
/etc/init.d/cron
/etc/pam.d/cron
/etc/default/cron
/etc/cron.d/.placeholder
/etc/cron.monthly/.placeholder
